# .claude/scripts/schedule/fredis-auth-watchdog.service
[Unit]
Description=Fredis auth watchdog (Claude + Google credential probe)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Replace __REPO_ROOT__ with the absolute path to the repo checkout on this
# host (e.g. /root/claude-code-second-brain).
WorkingDirectory=__REPO_ROOT__/.claude/scripts
ExecStart=/bin/bash __REPO_ROOT__/.claude/scripts/run_auth_watchdog.sh
# systemd strips PATH to /usr/bin:/bin by default; uv + claude live in ~/.local/bin.
Environment="PATH=/root/.local/bin:/usr/local/bin:/usr/bin:/bin"
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
